# timelapse_photo.crontab
* 9-18 * * 1-5 /path/to/timelapse_photo >/dev/null 2>&1
